# root/etc/services.d/jenkins/run
#!/usr/bin/with-contenv bash

echo "Starting jenkins slave agent..."

exec \
    s6-setuidgid lineage /bin/bash -lc \
    "java -jar /usr/local/bin/slave.jar -jnlpUrl http://192.168.1.88:8080/computer/${NAME}/slave-agent.jnlp -secret $JENKINS_AUTH_TOKEN"

